# scripts/com.medpact.overnight-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.medpact.overnight-build</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/christopherwilliams/Projects/medpact-practice-intelligence/scripts/overnight-build.sh</string>
    </array>
    
    <!-- Run at 2:00 AM every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/christopherwilliams/Projects/medpact-practice-intelligence/logs/launchd-out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/christopherwilliams/Projects/medpact-practice-intelligence/logs/launchd-err.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/christopherwilliams/Projects/medpact-practice-intelligence</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <!-- Keep trying if it fails -->
    <key>KeepAlive</key>
    <false/>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
